# backend/xiaoyuweihan-backend.service
[Unit]
Description=小雨微寒个人网站后端服务
Documentation=https://github.com/XYWH/xiaoyuweihan
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=notify
User=www
Group=www
RuntimeDirectory=xiaoyuweihan-backend
WorkingDirectory=/www/wwwroot/xiaoyuweihan/backend
Environment=PATH=/usr/local/bin/python3.10:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/www/wwwroot/xiaoyuweihan/backend
ExecStart=/usr/local/bin/python3.10 -m gunicorn --config gunicorn.conf.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=always
RestartSec=10

# 安全配置
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/www/wwwroot/xiaoyuweihan/backend
ReadWritePaths=/tmp
ReadWritePaths=/var/log

# 资源限制
LimitNOFILE=65536
LimitNPROC=32768
LimitCORE=0

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xiaoyuweihan-backend

[Install]
WantedBy=multi-user.target