# backend/xiaoyuweihan-backend.service
# systemd服务配置文件 - 宝塔版本
# /etc/systemd/system/xiaoyuweihan-backend.service

[Unit]
Description=小雨微寒 Backend Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=notify
User=www
Group=www
WorkingDirectory=/www/wwwroot/xiaoyuweihan/backend
Environment="PATH=/www/wwwroot/xiaoyuweihan/backend/venv/bin"
ExecStart=/www/wwwroot/xiaoyuweihan/backend/venv/bin/gunicorn -c gunicorn.conf.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=5

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/www/wwwroot/xiaoyuweihan/backend/uploads /www/wwwroot/xiaoyuweihan/backend/logs

# 资源限制
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target